# files/udev/rules.d/50-loop.rules
KERNEL=="loop*", RUN+="/etc/udev/scripts/loop-setup"
